# ktweak.service
[Unit]
Description=KTweak

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=ktweak
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
